# glean/init/glean.openrc
#!/sbin/runscript
# Copyright 1999-2016 Gentoo Foundation
# $Id$

description="basic instance setup"

start() {
  /usr/bin/glean.sh
  eend 0
}

depend() {
  need root localmount
  before net
}
